# mdebian-tnt/DEBIAN/postinst
#! /bin/sh

# Halt on error
set -e

# Active the streaming service
systemctl enable tnt.service
systemctl start tnt.service
