# mdebian-tnt/DEBIAN/prerm
#! /bin/sh

# Halt on error
set -e

# Disable the streaming service
systemctl stop tnt.service
systemctl disable tnt.service
